# com.alamontagne.transcribewatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Unique identifier — must match the filename (without .plist) -->
    <key>Label</key>
    <string>com.alamontagne.transcribe-watcher</string>

    <!-- The command to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/alamontagne/whisperx-env/bin/python3</string>
        <string>/Users/alamontagne/Documents/Trancscribe/watcher.py</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/alamontagne/Documents/Trancscribe</string>

    <!-- Environment variables available to the process.
         Add your real HF token value here so the watcher always has it,
         even if ~/.zshrc hasn't been sourced. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/alamontagne</string>
        <!-- IMPORTANT: replace the value below with your actual HF token -->
        <key>HF_TOKEN</key>
        <string>hf_enter-key-here</string>
    </dict>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Log files — tailed with: tail -f ~/Library/Logs/transcribe-watcher.out -->
    <key>StandardOutPath</key>
    <string>/Users/alamontagne/Library/Logs/transcribe-watcher.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/alamontagne/Library/Logs/transcribe-watcher.err</string>

    <!-- Throttle restarts — don't spin-loop if it keeps crashing -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

</dict>
</plist>
